# templates/pgbouncer.socket.j2
[Unit]
Description=sockets for PgBouncer

[Socket]
ListenStream={{ pgbouncer_listen_port }}
ListenStream=%i
ListenStream=/tmp/.s.PGSQL.%i

ReusePort=true

[Install]
WantedBy=sockets.target